# service/correlator/run
#!/bin/bash
cd ../..
# PYTHONPATH in cwd is for settings.py
PYTHONPATH=. exec ./bin/correlator
